# browsercontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:xvfb]
command=Xvfb :2 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
user=devuser
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display :2 -forever -shared -nopw -rfbport 5903 -noxdamage
user=devuser
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-server]
command=node /opt/browsercontainer/server.js
user=devuser
autostart=true
autorestart=true
priority=30
startsecs=5
directory=/opt/browsercontainer
environment=DISPLAY=":2",SIDECAR_TRANSPORT="sse",SIDECAR_PORT="8931"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
