# xr-runtime/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# One-shot root init: named volumes (cargo registry, gdext target) mount
# root-owned, so chown them to devuser before any build/run drops privileges.
[program:init-perms]
command=bash -c "chown -R devuser:devuser /home/devuser/.cargo /workspace/xr-client/rust/target 2>/dev/null; exit 0"
user=root
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xvfb]
command=Xvfb :3 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
user=devuser
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display :3 -forever -shared -nopw -rfbport 5904 -noxdamage
user=devuser
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:monado]
command=/opt/xr-runtime/launch-monado.sh
user=devuser
autostart=true
autorestart=true
priority=25
startsecs=5
environment=DISPLAY=":3",HOME="/home/devuser",XDG_RUNTIME_DIR="/tmp/xdg-runtime",XR_INPUT_DRIVER="%(ENV_XR_INPUT_DRIVER)s",XR_RUNTIME_HEADLESS="%(ENV_XR_RUNTIME_HEADLESS)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:godot]
command=/opt/xr-runtime/launch-godot.sh
user=devuser
autostart=true
autorestart=true
priority=30
# Generous startsecs: first boot compiles the gdext cdylib before Godot runs.
startsecs=10
startretries=10
environment=DISPLAY=":3",HOME="/home/devuser",XDG_RUNTIME_DIR="/tmp/xdg-runtime",XR_RUNTIME_JSON="%(ENV_XR_RUNTIME_JSON)s",XR_GODOT_SCENE="%(ENV_XR_GODOT_SCENE)s",XR_BACKEND_WS="%(ENV_XR_BACKEND_WS)s",LIBCLANG_PATH="/usr/lib",CARGO_HOME="/home/devuser/.cargo"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
