# config/supervisord.solid-pod.conf
; Optional: run solid-pod-rs-server directly inside the agentbox container
; rather than as a Docker sidecar. Copy the compiled binary to /usr/local/bin/
; then include this fragment in your supervisord config.
;
; Prerequisites:
;   - solid-pod-rs-server binary at /usr/local/bin/solid-pod-rs-server
;     (build with: cargo build --release --features git -p solid-pod-rs-server)
;   - /var/lib/solid/pods directory writable by devuser (uid 1000)

[program:solid-pod-server]
command=/usr/local/bin/solid-pod-rs-server
    --host 0.0.0.0
    --port 8410
user=devuser
environment=
    RUST_LOG="info",
    JSS_BASE_URL="%(ENV_SOLID_POD_PUBLIC_URL)s",
    JSS_STORAGE_TYPE="fs",
    JSS_STORAGE_ROOT="/var/lib/solid/pods",
    SOLID_ALLOWED_ORIGINS="%(ENV_SOLID_ALLOWED_ORIGINS)s",
    SOLID_ADMIN_KEY="%(ENV_SOLID_ADMIN_KEY)s"
autostart=false
autorestart=true
startretries=3
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/solid-pod-server.log
stderr_logfile=/var/log/supervisor/solid-pod-server.log
stdout_logfile_maxbytes=10MB
redirect_stderr=false
