# config/supervisord.conf
# Agentbox Minimal Supervisord Configuration
# ARM64 + x86_64 compatible, headless deployment

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0766

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================================================
# SSH Server (Port 22)
# ============================================================================

[program:sshd]
command=/usr/sbin/sshd -D -e
user=root
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.error.log

# ============================================================================
# Management API (Port 9090)
# ============================================================================

[program:management-api]
command=/usr/local/bin/node /opt/management-api/server.js
directory=/opt/management-api
user=devuser
environment=HOME="/home/devuser",PORT="9090",NODE_ENV="production"
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/management-api.log
stderr_logfile=/var/log/management-api.error.log

# ============================================================================
# Claude Z.AI Service (Port 9600 - Internal Only)
# Cost-effective Claude API proxy
# ============================================================================

[program:claude-zai]
command=/usr/local/bin/node /opt/claude-zai/server.js
directory=/opt/claude-zai
user=zai-user
environment=HOME="/home/zai-user",PORT="9600",NODE_ENV="production"
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/claude-zai.log
stderr_logfile=/var/log/claude-zai.error.log

# ============================================================================
# MCP Gateway (TCP + WebSocket)
# Claude Flow V3 MCP Infrastructure
# ============================================================================

[program:mcp-gateway]
command=/usr/local/bin/node /opt/mcp/servers/mcp-gateway.js
directory=/opt/mcp/servers
user=devuser
environment=HOME="/home/devuser",MCP_TCP_PORT="9500",MCP_WS_PORT="3002",LOG_LEVEL="info"
autostart=true
autorestart=true
priority=150
stdout_logfile=/var/log/mcp-gateway.log
stderr_logfile=/var/log/mcp-gateway.error.log

# ============================================================================
# code-server (Port 8080) - Optional web IDE
# Start manually: supervisorctl start code-server
# ============================================================================

[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth none /workspace
user=devuser
environment=HOME="/home/devuser",USER="devuser"
autostart=false
autorestart=true
priority=300
stdout_logfile=/var/log/code-server.log
stderr_logfile=/var/log/code-server.error.log

# ============================================================================
# Optional MCP Servers (on-demand)
# Start with: supervisorctl start <service>
# ============================================================================

[program:web-summary-mcp]
command=/usr/local/bin/node /opt/skills/web-summary/mcp-server/server.js
directory=/opt/skills/web-summary/mcp-server
user=devuser
environment=HOME="/home/devuser",ZAI_CONTAINER_URL="http://localhost:9600"
autostart=false
autorestart=true
priority=400
stdout_logfile=/var/log/web-summary-mcp.log
stderr_logfile=/var/log/web-summary-mcp.error.log

[program:imagemagick-mcp]
command=/usr/bin/python3 -u /opt/skills/imagemagick/tools/imagemagick_mcp.py
directory=/opt/skills/imagemagick/tools
user=devuser
environment=HOME="/home/devuser"
autostart=false
autorestart=true
priority=401
stdout_logfile=/var/log/imagemagick-mcp.log
stderr_logfile=/var/log/imagemagick-mcp.error.log

[program:playwright-mcp]
command=/usr/local/bin/node /opt/skills/playwright/tools/playwright_mcp.js
directory=/opt/skills/playwright/tools
user=devuser
environment=HOME="/home/devuser",PLAYWRIGHT_BROWSERS_PATH="/opt/playwright-browsers"
autostart=false
autorestart=true
priority=402
stdout_logfile=/var/log/playwright-mcp.log
stderr_logfile=/var/log/playwright-mcp.error.log
